# cron/tootbot.cron
*/30 * * * * python3 tootbot.py --operation toot --username <username> --instance <instance> --source <source>
