# cron/tootbot.cron
*/30 * * * * python3 /usr/src/app/tootbot.py --operation toot --username <username> --instance <instance> --source <source>
